# linux/ubuntu 18.04/etc/profile.d/android.sh
#!/bin/bash
export ANDROID_HOME='~/Applications/android-sdk'
