# linux/ubuntu 18.04/etc/profile.d/composer.sh
#!/bin/bash
export COMPOSER_HOME='~/.composer'
