# linux/ubuntu 18.04/etc/profile.d/java.sh
#!/bin/bash
export JDK8_DIR='~/.java/jdk1.8.0_191'
export JDK9_DIR='~/.java/jdk-9'
export JAVA_HOME=$JDK8_DIR
