# linux/ubuntu 18.04/etc/profile.d/oracle.sh
#!/bin/bash
export ORACLE_HOME=~/Applications/oracle/instantclient_12_2
